# claude-bot.service
[Unit]
Description=Claude Code Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=claude
Group=claude
WorkingDirectory=/home/claude/claude-bot
ExecStart=/home/claude/.local/bin/uv run --project /home/claude/claude-bot claude-bot
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-bot

# Безопасность
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/home/claude
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
