# claude-bot.service
[Unit]
Description=Claude Code Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bgs
Group=bgs
WorkingDirectory=/opt/claude_bot
ExecStart=/home/bgs/.local/bin/uv run --project /opt/claude_bot claude-bot
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-bot

# Безопасность
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/opt/claude_bot /home/bgs
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
